# etc/rc.d/amd
#!/bin/sh
#
# $OpenBSD: amd,v 1.7 2015/07/15 13:50:54 ajacoutot Exp $

daemon="/usr/sbin/amd"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"
rc_reload=NO
rc_stop=NO

rc_pre() {
	[ -e ${amd_master} ] || return 1
	daemon_flags="${daemon_flags} $(print -rn -- $(< ${amd_master}))"
}

rc_start() {
	${rcexec} "cd /etc/amd; ${daemon} ${daemon_flags}"
}

rc_cmd $1
