# etc/rc.d/amd
#!/bin/ksh
#
# $OpenBSD: amd,v 1.10 2022/05/21 12:42:21 ajacoutot Exp $

daemon="/usr/sbin/amd"

. /etc/rc.d/rc.subr

pexp="${daemon}.*"
rc_reload=NO
rc_stop=NO

rc_pre() {
	[[ -e ${amd_master} ]] || return 1
	daemon_flags="${daemon_flags} $(print -rn -- $(< ${amd_master}))"
}

rc_start() {
	rc_exec "cd /etc/amd; ${daemon} ${daemon_flags}"
}

rc_cmd $1
