# etc/rc.d/aucat
#!/bin/sh
#
# $OpenBSD: aucat,v 1.2 2011/07/08 05:16:08 ajacoutot Exp $

daemon="/usr/bin/aucat"

. /etc/rc.d/rc.subr

pexp="${daemon} -l${daemon_flags:+ ${daemon_flags}}"
rc_reload=NO

rc_pre() {
	daemon_flags="-l ${daemon_flags}"
}

rc_cmd $1
