# etc/rc.d/bgpd
#!/bin/sh
#
# $OpenBSD: bgpd,v 1.2 2015/12/19 13:45:12 ajacoutot Exp $

daemon="/usr/sbin/bgpd"

. /etc/rc.d/rc.subr

pexp="bgpd: parent.*"

# child will not return a config parsing error to the parent
rc_pre() {
	${daemon} -n ${daemon_flags}
}

rc_cmd $1
