# etc/rc.d/bootparamd
#!/bin/ksh
#
# $OpenBSD: bootparamd,v 1.3 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/rpc.bootparamd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	[ -s /etc/bootparams ]
}

rc_cmd $1
