# etc/rc.d/btd
#!/bin/sh
#
# $OpenBSD: btd,v 1.2 2011/07/08 04:29:54 ajacoutot Exp $

daemon="/usr/sbin/btd"

. /etc/rc.d/rc.subr

pexp="btd: \[priv\]"
rc_reload=NO

rc_cmd $1
