# etc/rc.d/cron
#!/bin/sh
#
# $OpenBSD: cron,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/cron"

. /etc/rc.d/rc.subr

rc_cmd $1
