# etc/rc.d/dhcpleased
#!/bin/ksh
#
# $OpenBSD: dhcpleased,v 1.1 2021/02/26 17:18:41 florian Exp $

daemon="/sbin/dhcpleased"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
