# etc/rc.d/dhcpd
#!/bin/sh
#
# $OpenBSD: dhcpd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/dhcpd"

. /etc/rc.d/rc.subr

rc_pre() {
	touch /var/db/dhcpd.leases
}

rc_reload=NO

rc_cmd $1
