# etc/rc.d/dvmrpd
#!/bin/sh
#
# $OpenBSD: dvmrpd,v 1.3 2016/02/07 21:00:16 sthen Exp $

daemon="/usr/sbin/dvmrpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
