# etc/rc.d/eigrpd
#!/bin/sh

daemon="/usr/sbin/eigrpd"

. /etc/rc.d/rc.subr

pexp="eigrpd: parent.*"

rc_cmd $1
