# etc/rc.d/eigrpd
#!/bin/sh
#
# $OpenBSD: eigrpd,v 1.2 2015/10/21 11:28:02 gsoares Exp $

daemon="/usr/sbin/eigrpd"

. /etc/rc.d/rc.subr

pexp="eigrpd: parent.*"

rc_cmd $1
