# etc/rc.d/ftpd
#!/bin/sh
#
# $OpenBSD: ftpd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/libexec/ftpd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
