# etc/rc.d/ftpd
#!/bin/sh
#
# $OpenBSD: ftpd,v 1.3 2011/12/21 09:26:29 ajacoutot Exp $

daemon="/usr/libexec/ftpd -D"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
