# etc/rc.d/hostapd
#!/bin/sh
#
# $OpenBSD: hostapd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/hostapd"

. /etc/rc.d/rc.subr

rc_cmd $1
