# etc/rc.d/httpd
#!/bin/ksh
#
# $OpenBSD: httpd,v 1.6 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/httpd"

. /etc/rc.d/rc.subr

# child will not return a config parsing error to the parent
rc_pre() {
	${daemon} -n ${daemon_flags}
}

rc_cmd $1
