# etc/rc.d/identd
#!/bin/sh
#
# $OpenBSD: identd,v 1.5 2013/08/07 22:07:35 dlg Exp $

daemon="/usr/sbin/identd"
daemon_flags=""

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
