# etc/rc.d/identd
#!/bin/sh
#
# $OpenBSD: identd,v 1.6 2013/08/08 15:41:28 ajacoutot Exp $

daemon="/usr/sbin/identd"
daemon_flags="-e"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
