# etc/rc.d/iked
#!/bin/ksh
#
# $OpenBSD: iked,v 1.6 2018/01/11 22:44:44 rpe Exp $

daemon="/sbin/iked"

. /etc/rc.d/rc.subr

pexp="iked: parent.*"

# Child will not return a config parsing error to the parent.
rc_pre() {
	[[ ${sasyncd_flags} != NO ]] &&
		daemon_flags="-S ${daemon_flags}"

	${daemon} -n ${daemon_flags}
}

rc_cmd $1
