# etc/rc.d/iked
#!/bin/sh
#
# $OpenBSD: iked,v 1.2 2015/12/19 10:09:04 ajacoutot Exp $

daemon="/sbin/iked"

. /etc/rc.d/rc.subr

pexp="iked: parent.*"

rc_pre() {
	[ X"${sasyncd_flags}" != X"NO" ] && \
		daemon_flags="-S ${daemon_flags}"
	return 0
}

rc_cmd $1
