# etc/rc.d/iked
#!/bin/ksh
#
# $OpenBSD: iked,v 1.4 2018/01/11 19:52:12 rpe Exp $

daemon="/sbin/iked"

. /etc/rc.d/rc.subr

pexp="iked: parent.*"

rc_pre() {
	[ X"${sasyncd_flags}" != X"NO" ] && \
		daemon_flags="-S ${daemon_flags}"
	#return 0
	# child will not return a config parsing error to the parent
	${daemon} -n ${daemon_flags}
}

rc_cmd $1
