# etc/rc.d/ipropd_master
#!/bin/sh
#
# $OpenBSD: ipropd_master,v 1.1 2013/08/09 16:24:54 ajacoutot Exp $

daemon="/usr/libexec/ipropd-master --detach"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
