# etc/rc.d/isakmpd
#!/bin/ksh
#
# $OpenBSD: isakmpd,v 1.3 2018/01/11 21:09:26 rpe Exp $

daemon="/sbin/isakmpd"

. /etc/rc.d/rc.subr

pexp="isakmpd: monitor \[priv\]"

rc_pre() {
	[[ ${sasyncd_flags} != NO ]] &&
		daemon_flags="-S ${daemon_flags}"
	return 0
}

rc_cmd $1
