# etc/rc.d/isakmpd
#!/bin/ksh
#
# $OpenBSD: isakmpd,v 1.2 2018/01/11 19:52:12 rpe Exp $

daemon="/sbin/isakmpd"

. /etc/rc.d/rc.subr

pexp="isakmpd: monitor \[priv\]"

rc_pre() {
	[ X"${sasyncd_flags}" != X"NO" ] && \
		daemon_flags="-S ${daemon_flags}"
	return 0
}

rc_cmd $1
