# etc/rc.d/iscsid
#!/bin/ksh
#
# $OpenBSD: iscsid,v 1.2 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/iscsid"
control="/usr/sbin/iscsictl"

. /etc/rc.d/rc.subr

rc_reload=NO
rc_stop=NO

rc_start() {
	${rcexec} "${daemon} ${daemon_flags} && ${control} reload"
}

rc_cmd $1
