# etc/rc.d/kdc
#!/bin/sh
#
# $OpenBSD: kdc,v 1.2 2013/06/17 19:23:24 ajacoutot Exp $

daemon="/usr/libexec/kdc --detach"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
