# etc/rc.d/ldapd
#!/bin/sh
#
# $OpenBSD: ldapd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/ldapd"

. /etc/rc.d/rc.subr

rc_cmd $1
