# etc/rc.d/lockd
#!/bin/sh
#
# $OpenBSD: lockd,v 1.1 2011/07/08 00:54:04 ajacoutot Exp $

daemon="/usr/sbin/rpc.lockd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
