# etc/rc.d/lpd
#!/bin/sh
#
# $OpenBSD: lpd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/lpd"

. /etc/rc.d/rc.subr

rc_cmd $1
