# etc/rc.d/mopd
#!/bin/sh
#
# $OpenBSD: mopd,v 1.2 2011/07/08 02:15:34 robert Exp $

daemon="/usr/sbin/mopd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	[ -d /tftpboot/mop ]
}

rc_cmd $1
