# etc/rc.d/mountd
#!/bin/sh
#
# $OpenBSD: mountd,v 1.3 2012/08/14 07:43:57 ajacoutot Exp $

daemon="/sbin/mountd"

. /etc/rc.d/rc.subr

rc_pre() {
	[ -s /etc/exports ] && grep -qv '^#' /etc/exports && \
		rm -f /var/db/mountdtab && \
		echo -n > /var/db/mountdtab
}

rc_cmd $1
