# etc/rc.d/named
#!/bin/sh
#
# $OpenBSD: named,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/named"

. /etc/rc.d/rc.subr

pexp="named: \[priv\]"

rc_cmd $1
