# etc/rc.d/nfsd
#!/bin/sh
#
# $OpenBSD: nfsd,v 1.2 2011/08/02 11:50:39 robert Exp $

daemon="/sbin/nfsd"

. /etc/rc.d/rc.subr

pexp="(${daemon}*|nfsd: (master|server)$)"

rc_reload=NO

rc_pre() {
	[ -s /etc/exports -a `sed -e '/^#/d' < /etc/exports | wc -l` -ne 0 ] &&
		rm -f /var/db/mountdtab && \
		echo -n > /var/db/mountdtab
}

rc_cmd $1
