# etc/rc.d/npppd
#!/bin/sh
#
# $OpenBSD: npppd,v 1.1 2012/09/20 12:51:43 yasuoka Exp $

daemon="/usr/sbin/npppd"

. /etc/rc.d/rc.subr

pexp="npppd: main"

rc_cmd $1
