# etc/rc.d/npppd
#!/bin/sh
#
# $OpenBSD: npppd,v 1.2 2016/02/02 17:51:11 sthen Exp $

daemon="/usr/sbin/npppd"

. /etc/rc.d/rc.subr

rc_cmd $1
