# etc/rc.d/nsd
#!/bin/sh
#
# $OpenBSD: nsd,v 1.2 2011/12/13 13:45:38 ajacoutot Exp $

daemon="/usr/sbin/nsd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
