# etc/rc.d/ospf6d
#!/bin/ksh
#
# $OpenBSD: ospf6d,v 1.4 2026/03/17 14:27:25 sthen Exp $

daemon="/usr/sbin/ospf6d"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
