# etc/rc.d/pflogd
#!/bin/sh
#
# $OpenBSD: pflogd,v 1.2 2011/07/08 21:42:52 robert Exp $

daemon="/sbin/pflogd"

. /etc/rc.d/rc.subr

pexp="pflogd: \[priv\]"

rc_pre() {
	if pfctl -si | grep -q Enabled; then
		ifconfig pflog0 create
		if ifconfig pflog0; then
			ifconfig pflog0 up
		else
			return 1
		fi
	else
		return 1
	fi
}

rc_cmd $1
