# etc/rc.d/popa3d
#!/bin/sh
#
# $OpenBSD: popa3d,v 1.2 2011/09/19 20:42:24 ajacoutot Exp $

daemon="/usr/sbin/popa3d"
daemon_flags="-D"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
