# etc/rc.d/popa3d
#!/bin/sh
#
# $OpenBSD: popa3d,v 1.3 2011/12/21 09:26:29 ajacoutot Exp $

daemon="/usr/sbin/popa3d -D"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
