# etc/rc.d/relayd
#!/bin/ksh
#
# $OpenBSD: relayd,v 1.7 2022/05/21 12:42:21 ajacoutot Exp $

daemon="/usr/sbin/relayd"

. /etc/rc.d/rc.subr

# Child will not return a config parsing error to the parent.
rc_pre() {
	# use rc_exec here since daemon_flags may contain arguments with spaces
	rc_exec "${daemon} -n ${daemon_flags}"
}

rc_cmd $1
