# etc/rc.d/relayd
#!/bin/sh
#
# $OpenBSD: relayd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/relayd"

. /etc/rc.d/rc.subr

pexp="relayd: parent.*"

rc_cmd $1
