# etc/rc.d/route6d
#!/bin/sh
#
# $OpenBSD: route6d,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/route6d"

. /etc/rc.d/rc.subr

rc_cmd $1
