# etc/rc.d/rtadvd
#!/bin/sh
#
# $OpenBSD: rtadvd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/rtadvd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
