# etc/rc.d/rtsold
#!/bin/sh
#
# $OpenBSD: rtsold,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/rtsold"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
