# etc/rc.d/sasyncd
#!/bin/sh
#
# $OpenBSD: sasyncd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/sasyncd"

. /etc/rc.d/rc.subr

pexp="sasyncd: \[priv\]"

rc_cmd $1
