# etc/rc.d/sendmail
#!/bin/sh
#
# $OpenBSD: sendmail,v 1.5 2012/01/20 22:27:30 sthen Exp $

daemon="/usr/sbin/sendmail"

. /etc/rc.d/rc.subr

pexp="(sendmail: (accepting|rejecting).*|${daemon}* -(q[0-9]|bd)*)"

rc_bg=YES

rc_cmd $1
