# etc/rc.d/sendmail
#!/bin/sh
#
# $OpenBSD: sendmail,v 1.6 2012/07/10 15:00:02 halex Exp $

daemon="/usr/sbin/sendmail"

. /etc/rc.d/rc.subr

pexp="(sendmail: (accepting|rejecting).*|${daemon}.* -(q[0-9]|bd)*)"

rc_bg=YES

rc_cmd $1
