# etc/rc.d/sendmail
#!/bin/sh
#
# $OpenBSD: sendmail,v 1.4 2011/07/12 05:40:55 robert Exp $

daemon="/usr/sbin/sendmail"

. /etc/rc.d/rc.subr

pexp="(sendmail: accepting.*|${daemon}* -(q[0-9]|bd)*)"

rc_bg=YES

rc_cmd $1
