# etc/rc.d/sndiod
#!/bin/ksh
#
# $OpenBSD: sndiod,v 1.5 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/bin/sndiod"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
