# etc/rc.d/sndiod
#!/bin/sh
#
# $OpenBSD: sndiod,v 1.4 2016/04/28 09:15:16 ajacoutot Exp $

daemon="/usr/bin/sndiod"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
