# etc/rc.d/spamd
#!/bin/sh
#
# $OpenBSD: spamd,v 1.3 2013/09/13 14:50:56 okan Exp $

daemon="/usr/libexec/spamd"

. /etc/rc.d/rc.subr

pexp="spamd: \[priv\]"
rc_reload=NO

rc_pre() {
	[ X"${spamd_black}" != X"NO" ] && \
		daemon_flags="-b ${daemon_flags}"
	return 0
}

rc_cmd $1
