# etc/rc.d/spamlogd
#!/bin/ksh
#
# $OpenBSD: spamlogd,v 1.4 2018/01/11 21:09:26 rpe Exp $

daemon="/usr/libexec/spamlogd"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	[[ ${spamd_flags} != NO && ${spamd_black} == NO ]]
	if pfctl -si | grep -q Enabled; then
		ifconfig pflog0 create
		if ifconfig pflog0; then
			ifconfig pflog0 up
		else
			return 1
		fi
	else
		return 1
	fi
}

rc_cmd $1
