# etc/rc.d/sshd
#!/bin/ksh
#
# $OpenBSD: sshd,v 1.4 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/sshd"

. /etc/rc.d/rc.subr

rc_reload() {
	${daemon} ${daemon_flags} -t && pkill -HUP -xf "${pexp}"
}

rc_cmd $1
