# etc/rc.d/sshd
#!/bin/ksh
#
# $OpenBSD: sshd,v 1.5 2020/01/22 13:14:51 sthen Exp $

daemon="/usr/sbin/sshd"

. /etc/rc.d/rc.subr

pexp="sshd: \[listener\].*"

rc_reload() {
	${daemon} ${daemon_flags} -t && pkill -HUP -xf "${pexp}"
}

rc_cmd $1
