# etc/rc.d/sshd
#!/bin/sh
#
# $OpenBSD: sshd,v 1.1 2011/07/06 18:55:36 robert Exp $

daemon="/usr/sbin/sshd"

. /etc/rc.d/rc.subr

rc_cmd $1
