# etc/rc.d/switchd
#!/bin/ksh
#
# $OpenBSD: switchd,v 1.3 2018/01/11 22:44:44 rpe Exp $

daemon="/usr/sbin/switchd"

. /etc/rc.d/rc.subr

# Child will not return a config loading error to the parent.
rc_pre() {
	${daemon} -n ${daemon_flags}
}

rc_cmd $1
