# etc/rc.d/tftpproxy
#!/bin/sh
#
# $OpenBSD: tftpproxy,v 1.2 2015/07/15 13:50:54 ajacoutot Exp $

daemon="/usr/sbin/tftp-proxy"

. /etc/rc.d/rc.subr

pexp="tftp-proxy: privproc"
rc_reload=NO

rc_cmd $1
