# etc/rc.d/tftpd
#!/bin/sh
#
# $OpenBSD: tftpd,v 1.2 2012/07/07 09:35:36 ajacoutot Exp $

daemon="/usr/sbin/tftpd"
daemon_flags="/tftpboot"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
