# etc/rc.d/unbound
#!/bin/ksh
#
# $OpenBSD: unbound,v 1.8 2024/09/23 20:44:24 kn Exp $

daemon="/usr/sbin/unbound"
daemon_flags="-c /var/unbound/etc/unbound.conf"

. /etc/rc.d/rc.subr

rc_pre() {
	if grep '^[[:space:]]*auto-trust-anchor-file:' \
            /var/unbound/etc/unbound.conf > /dev/null 2>&1; then
		/usr/sbin/unbound-anchor -v
	fi

	/usr/sbin/unbound-checkconf
}

rc_cmd $1
