# etc/rc.d/unbound
#!/bin/ksh
#
# $OpenBSD: unbound,v 1.7 2022/05/21 10:10:45 ajacoutot Exp $

daemon="/usr/sbin/unbound"
daemon_flags="-c /var/unbound/etc/unbound.conf"

. /etc/rc.d/rc.subr

rc_pre() {
	if grep '^[[:space:]]*auto-trust-anchor-file:' \
            /var/unbound/etc/unbound.conf > /dev/null 2>&1; then
		/usr/sbin/unbound-anchor -v
	fi
	/usr/sbin/unbound-checkconf || return 1
}

rc_cmd $1
