# etc/rc.d/unwind
#!/bin/ksh
#
# $OpenBSD: unwind,v 1.1 2019/01/26 10:59:07 florian Exp $

daemon="/sbin/unwind"

. /etc/rc.d/rc.subr

rc_pre() {
	/usr/sbin/unbound-anchor -a /etc/unwind/trustanchor/root.key || true
}

rc_cmd $1
