# etc/rc.d/ypbind
#!/bin/sh
#
# $OpenBSD: ypbind,v 1.4 2011/10/07 05:54:48 ajacoutot Exp $

daemon="/usr/sbin/ypbind"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_pre() {
	[ X"`domainname`" != X"" ]
}

rc_cmd $1
