# etc/rc.d/ypbind
#!/bin/sh
#
# $OpenBSD: ypbind,v 1.6 2017/05/27 13:51:52 ajacoutot Exp $

daemon="/usr/sbin/ypbind"

. /etc/rc.d/rc.subr

rc_reload=NO

rc_pre() {
	[ X"`domainname`" != X"" ]
}

rc_post() {
	rm -f /var/run/ypbind.lock
}

rc_cmd $1
