# etc/rc.d/ypserv
#!/bin/ksh
#
# $OpenBSD: ypserv,v 1.5 2018/01/11 19:52:12 rpe Exp $

daemon="/usr/sbin/ypserv"

. /etc/rc.d/rc.subr

rc_pre() {
	[ X"`domainname`" != X"" -a -d /var/yp/`domainname` ]
}

rc_cmd $1
