# etc/rc.d/ypserv
#!/bin/sh
#
# $OpenBSD: ypserv,v 1.4 2011/10/06 06:47:50 ajacoutot Exp $

daemon="/usr/sbin/ypserv"

. /etc/rc.d/rc.subr

rc_pre() {
	[ X"`domainname`" != X"" -a -d /var/yp/`domainname` ]
}

rc_cmd $1
