# etc/rc.local
#	$OpenBSD: rc.local,v 1.42 2011/03/17 16:43:51 robert Exp $

# Site-specific startup actions, daemons, and other things which
# can be done AFTER your system goes into securemode.  For actions
# which should be done BEFORE your system has gone into securemode
# please see /etc/rc.securelevel.

echo -n 'starting local daemons:'

for _r in $rc_scripts; do
	[ -x /etc/rc.d/${_r} ] && /etc/rc.d/${_r} start
done

# Add your local startup actions here.

echo '.'
