# etc/rc.local
#	$OpenBSD: rc.local,v 1.41 2010/11/05 10:03:00 ajacoutot Exp $

# Site-specific startup actions, daemons, and other things which
# can be done AFTER your system goes into securemode.  For actions
# which should be done BEFORE your system has gone into securemode
# please see /etc/rc.securelevel.

echo -n 'starting local daemons:'

for _r in $rc_scripts; do
	[ -x /etc/rc.d/${_r} ] && echo -n " ${_r}" && /etc/rc.d/${_r} start
done

# Add your local startup actions here.

echo '.'
